# gygax/s6/services.d/gygax/run
#!/usr/bin/with-contenv sh

cd /app
exec python server.py