# firehose/s6/services.d/slack-firehose/run
#!/usr/bin/with-contenv sh

cd /app
exec python reader.py