# auth/s6/services.d/auth/run
#!/usr/bin/with-contenv sh

cd /app
exec gunicorn --bind "0.0.0.0:${PORT}" --worker-class eventlet --workers 2 server:app
